# seed-knowledge-manager/agent/systemd/km-poll.timer
[Unit]
Description=Knowledge Manager — poll mentions every 60s

[Timer]
OnBootSec=90
OnUnitActiveSec=15
AccuracySec=5s
Persistent=true
Unit=km-poll.service

[Install]
WantedBy=timers.target
