# seed-knowledge-manager/agent/systemd/km-poll.service
[Unit]
Description=Knowledge Manager — poll mentions
After=default.target

[Service]
Type=oneshot
EnvironmentFile=/home/km/.nanobot/secrets.env
Environment=PATH=/home/km/.local/bin:/usr/bin:/usr/local/bin
Environment=KM_TRIGGER=poll-cli
Environment=SEED_CLI_PATH=/home/km/.local/bin/seed-cli
Environment=KM_STATE_DIR=/home/km/km-state
Environment=KM_LOGS_DIR=/home/km/km-logs
# Standalone polling driver. Bypasses nanobot — DeepSeek call happens
# inline, no LLM tool-orchestration. Whole loop is deterministic except
# for the single chat-completion that drafts the reply text.
ExecStart=/usr/bin/node /home/km/km-agent/mcp/seed-cli-mcp/dist/poll-cli.js
# Timeout per run. Plenty of headroom for one DeepSeek call + a few seed-cli writes.
TimeoutStartSec=180

[Install]
WantedBy=default.target
